# configs/stellarrobot_gunicorn.conf
# Supervisor config file
# contents of /etc/supervisor/conf.d/stellarrobot_gunicorn.conf
[program:stellarrobot_gunicorn]

# the command to start Gunicorn
command = /home/webuser/www/stellarrobot/venv/bin/python /home/webuser/www/stellarrobot/venv/bin/gunicorn -b unix:/tmp/stellarrobot.sock -w 3 -k gevent app:app

# the directory to be in while running this
directory = /home/webuser/www/stellarrobot/stellarrobot

# the user to run this service as
user = webuser

# start this at boot, and restart it if it fails
autostart = true
autorestart = true

stdout_logfile = /home/webuser/www/stellarrobot/logs/gunicorn.log
stderr_logfile = /home/webuser/www/stellarrobot/logs/gunicorn_err.log
